# com.resume-resume.site.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.resume-resume.site</string>
    <key>ProgramArguments</key>
    <array>
        <string>/Users/dshanklinbv/.pyenv/versions/3.12.7/bin/resume-site</string>
        <string>--host</string>
        <string>127.0.0.1</string>
        <string>--port</string>
        <string>8765</string>
    </array>
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <true/>
    <key>WorkingDirectory</key>
    <string>/Users/dshanklinbv</string>
    <key>StandardOutPath</key>
    <string>/Users/dshanklinbv/.claude/resume-site-stdout.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/dshanklinbv/.claude/resume-site-stderr.log</string>
    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/Users/dshanklinbv/.pyenv/versions/3.12.7/bin:/Users/dshanklinbv/.local/bin:/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>
        <key>HOME</key>
        <string>/Users/dshanklinbv</string>
    </dict>
  <key>ThrottleInterval</key>
  <integer>5</integer>
</dict>
</plist>
